# deploy/tell_nte.service
[Unit]
Description=TELL NTE – Textile Ecosystem Mapping (tell.newtexeco.nl)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/tell/app

EnvironmentFile=/home/tell/app/.env

ExecStart=/home/tell/app/venv/bin/gunicorn \
    --config /home/tell/app/deploy/gunicorn_conf.py \
    wsgi:application

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
