# deploy/tell-form.service
[Unit]
Description=TELL – Textile Ecosystem Contribution Form
After=network.target

[Service]
User=tell
Group=www-data
WorkingDirectory=/home/tell/app
ExecStart=/home/tell/app/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:8051 \
    --timeout 120 \
    contribution_form:server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
